# .config/systemd/user/backup.timer
[Unit]
Description=Runs every 1 time a day backup.sh

[Timer]
OnCalendar=*-*-* 20:00:00

[Install]
WantedBy=timers.target
